# dist/kbdtomidi.conf
# run the keyboard-to-midi daemon
#
# kbdtomidi takes exclusive control of an input device and generates midi
# events from the button presses instead of keyboard events.
# 
# This script it for whatever init system fedora used between sysv and systemd
# maybe it was upstart? I honestly don't remember. This probably isn't the file
# you want, take a look at the .conf file which is for systemd.


start on runlevel [345]

stop on runlevel [!5]

console output
respawn
respawn limit 100 120
exec /usr/local/share/gigrig/kbdtomidi/kbdtomidi /dev/input/by-id/usb-_USB_Keyboard-event-kbd virtual > /tmp/kbdtomidi.log 2>&1

pre-start script
    sleep 3	
end script

